# supervisord.conf
[program:api]
command=gunicorn -b 0.0.0.0:8060  'api.main:get_app()'
directory=/opt/api/
user=root
autostart=true
autorestart=false
redirect_stderr=true